# systemd/test_hf126.service
[Unit]
Description=Test daemon that does nothing special 

[Service]
ExecStart=@prefix@/bin/test_hf126
Type=notify
NotifyAccess=main
